# systemd/mounts/Laptop/var-home-ayush-Documents-University.automount
[Unit]
Description=Hardrive University Automount

[Automount]
Where=/var/home/ayush/Documents/University

[Install]
WantedBy=multi-user.target
