# systemd/mounts/Laptop/var-home-ayush-Videos-Movies.automount
[Unit]
Description=Hardrive Movies Automount

[Automount]
Where=/var/home/ayush/Videos/Movies

[Install]
WantedBy=multi-user.target
