# systemd/mounts/Laptop/var-home-ayush-Videos-Series.automount
[Unit]
Description=Hardrive Series Automount

[Automount]
Where=/var/home/ayush/Videos/Series

[Install]
WantedBy=multi-user.target
