# systemd/mounts/Laptop/var-home-ayush-Games.automount
[Unit]
Description=Hardrive Games Automount

[Automount]
Where=/var/home/ayush/Games

[Install]
WantedBy=multi-user.target
