# systemd/mounts/Laptop/var-home-ayush-Virtualization.automount
[Unit]
Description=SSD VM

[Automount]
Where=/var/home/ayush/Virtualization

[Install]
WantedBy=multi-user.target
