# systemd/mounts/var-home-ayush-Pictures-Harddrive.automount
[Unit]
Description=Hardrive Pictures Automount

[Automount]
Where=/var/home/ayush/Pictures/Harddrive

[Install]
WantedBy=multi-user.target
