# systemd/mounts/var-home-ayush-Documents-CollegeStuff.automount
[Unit]
Description=Hardrive College Stuff Automount

[Automount]
Where=/var/home/ayush/Documents/CollegeStuff

[Install]
WantedBy=multi-user.target
