# systemd/mounts/var-home-ayush-Videos-Lbry.automount
[Unit]
Description=Hardrive Videos-Lbry Automount

[Automount]
Where=/var/home/ayush/Videos/Lbry

[Install]
WantedBy=multi-user.target
