# systemd/mounts/var-home-ayush-Documents-Notes.automount
[Unit]
Description=Hardrive Documents-Notes Automount

[Automount]
Where=/var/home/ayush/Documents/Notes

[Install]
WantedBy=multi-user.target
