# systemd/mounts/var-home-ayush-Videos-Hollywood.automount
[Unit]
Description=Hardrive Videos-Hollywood Automount

[Automount]
Where=/var/home/ayush/Videos/Hollywood

[Install]
WantedBy=multi-user.target
