# systemd/mounts/Desktop/var-home-ayush-Pictures-Backup.automount
[Unit]
Description=Hardrive Pictures Backup

[Automount]
Where=/var/home/ayush/Pictures/Backup

[Install]
WantedBy=multi-user.target
