# systemd/mounts/Desktop/var-home-ayush-Documents-Books.automount
[Unit]
Description=Hardrive Books Automount

[Automount]
Where=/var/home/ayush/Documents/Books

[Install]
WantedBy=multi-user.target
