# systemd/mounts/Desktop/var-home-ayush-Documents-Programming.automount
[Unit]
Description=Hardrive Games Automount

[Automount]
Where=/var/home/ayush/Documents/Programming

[Install]
WantedBy=multi-user.target
