# systemd/mounts/Desktop/var-home-ayush-Videos-Anime.automount
[Unit]
Description=Hardrive Videos-Anime Automount

[Automount]
Where=/var/home/ayush/Videos/Anime

[Install]
WantedBy=multi-user.target
