# systemd/mounts/Desktop/var-home-ayush-Documents-Harddrive.automount
[Unit]
Description=Hardrive Documents-Harddrive Automount

[Automount]
Where=/var/home/ayush/Documents/Harddrive

[Install]
WantedBy=multi-user.target
