# systemd/mounts/Desktop/var-home-ayush-kde.automount
[Unit]
Description=Hardrive KDE Automount

[Automount]
Where=/var/home/ayush/kde

[Install]
WantedBy=multi-user.target
